# docker/supervisord.conf
; Runs the keygen backend API and static frontend as one foreground process
; tree (PID 1 = supervisord), so a single container needs no docker-compose.
; supervisord forwards SIGTERM from `docker stop` to both child processes.
;
; supervisord itself runs as root (no USER in the Dockerfile) because it
; needs to open /dev/stdout and /dev/stderr to wire up each child's logs —
; on Linux that's a root-owned, mode-0600 pipe that only root can (re)open
; by path, even though a non-root process can freely write to an *inherited*
; copy of the same fd. Each actual application process still drops to the
; unprivileged 'keygen' user via the per-program `user=` directive below.

[supervisord]
nodaemon=true
loglevel=info
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; --- backend: OpenSCAD-backed STL generation API (port 8080) ---------------
; Runs docker/serve_cors.py instead of bin/serve.py directly: it's a thin
; wrapper that adds Access-Control-Allow-Origin so the frontend (a different
; origin/port) is allowed to call it. See that file for details.
[program:keygen-api]
command=python3 docker/serve_cors.py
directory=/app
user=keygen
environment=HOME="/home/keygen"
autostart=true
autorestart=true
startsecs=1
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; --- frontend: static file server for web/ ----------------------------------
[program:keygen-web]
command=python3 -m http.server %(ENV_WEB_PORT)s --directory web
directory=/app
user=keygen
environment=HOME="/home/keygen"
autostart=true
autorestart=true
startsecs=1
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
